# deploy/obsidian-vault/systemd/obsidian-vault-healthcheck.service
[Unit]
Description=Obsidian Vault healthcheck
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
User=rodrigo
Group=rodrigo
WorkingDirectory=/data/apps/deploy/obsidian-vault
ExecStart=/data/apps/deploy/obsidian-vault/scripts/healthcheck.sh

[Install]
WantedBy=multi-user.target
